# gunicorn.conf
[program:gunicorn]
command=gunicorn run:app -b localhost:5000
directory=/blogsite
