# gunicorn.conf
[program:gunicorn]
command=gunicorn run:app -b localhost:8080
directory=/blogsite
